# moon-telegram-bot.service
[Unit]
Description=The Moon Ecosystem — MoonBot Telegram
Documentation=https://github.com/newjsouza/the-moon-ecosystem
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=johnathan
WorkingDirectory=/home/johnathan
ExecStart=/home/johnathan/moon_bot_start.sh
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moon-telegram-bot

[Install]
WantedBy=multi-user.target
