# the-moon.service
[Unit]
Description=The Moon Ecosystem — Autonomous AI Agent Loop
Documentation=https://github.com/newjsouza/the-moon-ecosystem
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=johnathan
WorkingDirectory=/home/johnathan
ExecStart=/home/johnathan/moon_start.sh
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=30
Restart=on-failure
RestartSec=15

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=the-moon

[Install]
WantedBy=multi-user.target